# src/main/resources/control/postinst
#!/bin/bash

# Launch at startup
sudo update-rc.d -f snowcamp-demo-0.3.0 defaults
