# src/main/resources/control/postinst
#!/bin/bash

sudo update-rc.d -f snowcamp-demo-0.2.2 defaults
# start services
# sudo service snowcamp-demo-0.2.0 start
